# com.datascore.agent-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.datascore.agent-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/datascore/agents/agent-sync-service.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/datascore/.claude/agent-sync.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/datascore/.claude/agent-sync.err</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/datascore/agents</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/datascore</string>
    </dict>
</dict>
</plist>
